# waykey.service
[Unit]
Description=WayKey Automation Engine (AutoHotkey alternative for Wayland)
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=uwsm app -- /usr/bin/node /opt/waykey/run.js
Restart=on-failure
RestartSec=5